# roles/zfs-monitoring/templates/systemd/zfs-scrub@.service.j2
#
# {{ ansible_managed }}
#

[Unit]
Description=zpool scrub on %i

[Service]
KillSignal=SIGINT
ExecStart=/sbin/zpool scrub -w %i

[Install]
WantedBy=multi-user.target
