# roles/zfs-monitoring/templates/systemd/zfs-trim@.service.j2
#
# {{ ansible_managed }}
#

[Unit]
Description=zpool trim on %i

[Service]
KillSignal=SIGINT
ExecStart=/usr/bin/zpool trim -w %i

[Install]
WantedBy=multi-user.target
